# crontab
* */1 * * * root /etc/cron.d/bitbucket-backup.sh >> /var/log/bitbucket-backup.log 2>&1
